# homestead.service
[Unit]
Description=Homestead — self-hosted dashboard
Documentation=https://github.com/your-username/homestead
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=homestead
Group=homestead

ExecStart=/usr/bin/homestead \
    --config /etc/homestead/config.toml \
    --port 8080

Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/homestead

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homestead

[Install]
WantedBy=multi-user.target
